# var/spool/cron/crontabs/www-data
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.WUPskY/crontab installed on Thu Dec  5 04:27:06 2019)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
MAILTO=""
* * * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento cron:run | grep -v "Ran jobs by schedule" >> /var/www/sites/canadasatellite.ca/httpdocs/var/log/magento.cron.log
* * * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/update/cron.php >> /var/www/sites/canadasatellite.ca/httpdocs/var/log/update.cron.log
* * * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento setup:cron:run >> /var/www/sites/canadasatellite.ca/httpdocs/var/log/setup.cron.log
#* * * * * touch /var/www/sites/canadasatellite.ca/httpdocs/var/test.txt
#0 3 * * *  /bin/sh /opt/magento_session_cleanup.sh > /dev/null 2>&1
#0 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento setup:static-content:deploy en_US en_CA fr_CA

#Magento live start
0 5 * * * /bin/sh /var/www/sites/canadasatellite.ca/cron.sh
2 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento currency:cron:save && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento indexer:reindex catalog_product_price && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento c:f
15 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5326 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5335 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5467 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5468 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5469 
17 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5470 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5471 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5473 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5474 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5475
19 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5476 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5479 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5642 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5658 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5667
21 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5733 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5745 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5750 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5764 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5908 
23 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5971 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5976 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5977 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5980 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=5994 
25 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6024 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6026 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6047 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6048 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6049
27 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6051 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6052 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6082 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6095 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6096 
29 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6097 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6098 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6102 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6104 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6107 
31 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=7548 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=7627 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=7870 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=7885 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=7891 
33 2 * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=8078 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=8113 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=8549 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=8579 && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento indexer:reindex catalog_product_price && /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento c:f
#* * * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento cron:run | grep -v Ran jobs by schedule >> /var/www/html/magento2/var/log/magento.cron.log
#* * * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/update/cron.php >> /var/www/html/magento2/var/log/update.cron.log
#* * * * * /usr/bin/php /var/www/sites/canadasatellite.ca/httpdocs/bin/magento setup:cron:run >> /var/www/html/magento2/var/log/setup.cron.log
#~ MAGENTO END

#Magento stage start
#2 2 * * * /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/bin/magento currency:cron:save && /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/bin/magento indexer:reindex catalog_product_price && /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/bin/magento c:f
13 6 * * * /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/bin/magento product:bundle:resave --Id=6024 && /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/bin/magento indexer:reindex catalog_product_price && /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/bin/magento c:f
#30 5 * * * /bin/sh /var/www/sites/staging.canadasatellite.ca/cron.sh
#* * * * * /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/bin/magento cron:run | grep -v Ran jobs by schedule >> /var/www/html/magento2/var/log/magento.cron.log
#* * * * * /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/update/cron.php >> /var/www/html/magento2/var/log/update.cron.log
#* * * * * /usr/bin/php /var/www/sites/staging.canadasatellite.ca/httpdocs/bin/magento setup:cron:run >> /var/www/html/magento2/var/log/setup.cron.log
#~ MAGENTO END


